# ansible/templates/rkhunter.service.j2
[Unit]
Description=Rootkit Scan

[Service]
Type=oneshot
ExecStart=/usr/bin/rkhunter --cronjob --propupd --report-warnings-only